# crontab
*/10 * * * * /home/mashod/Projects/CodeQuivers/portal-247/dbBackup_script-dev.sh